# _posts/2020-12-07-Backup Git Repositories/update-git-mirrors.timer
[Unit]
Description=Regularly refresh my Git mirrors

[Timer]
Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target
